# stackchan-relay.service
[Unit]
Description=StackChan relay + MCP server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/stackchan-relay
EnvironmentFile=/opt/stackchan-relay/.env
ExecStart=/opt/stackchan-relay/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8011
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
